# release/runeflow/DEBIAN/postinst
#!/bin/sh

systemctl enable runeflow.service
systemctl start runeflow.service
