# release/runeflow/DEBIAN/prerm
#!/bin/sh

systemctl stop runeflow.service
systemctl disable runeflow.service
