# release/runeflow/DEBIAN/prerm
#!/bin/sh

systemctl stop runeflow.service
systemctl disable runeflow.service

echo "args: $1"

if [ "$1" = "remove" ]; then
  echo "removing agent_id file"
  rm -rf /var/lib/runeflow
fi
